# supervisord.conf
[supervisord]
nodaemon=true

[program:naemon]
command=/usr/bin/naemon /etc/naemon/naemon.cfg
